# init.d/wg-broadcast
#!/sbin/openrc-run

supervisor=supervise-daemon
command="/etc/kislota-cloud/wg-broadcast.py"

depend() {
    need net
}
